# systemd/mtproxy-fork.service
[Unit]
Description=MTProxy Fork (hot secret reload)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mtproxy
Group=mtproxy
WorkingDirectory=/var/lib/mtproxy-fork
EnvironmentFile=-/etc/default/mtproxy-fork
ExecStart=/usr/local/bin/mtproxy-fork-run
Restart=on-failure
RestartSec=2
LimitNOFILE=262144
NoNewPrivileges=true
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
