# systemd/mtproxy-fork-expire-sync.service
[Unit]
Description=MTProxy Fork Expire Sync
After=mtproxy-fork.service
Wants=mtproxy-fork.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/proxyctl secret expire-disable
